# deploy/systemd/rustydb@.service
[Unit]
Description=RustyDB Instance %i
Documentation=https://github.com/harborgrid-justin/rusty-db
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Run as a dedicated user (create with: useradd -r -s /bin/false rustydb)
User=rustydb
Group=rustydb

# Instance root (Instance Layout Spec v1.0)
Environment=RUSTYDB_INSTANCE=%i
Environment=RUSTYDB_HOME=/var/lib/rustydb/instances/%i

# Main server process
ExecStart=/opt/rustydb/current/bin/rusty-db-server --home /var/lib/rustydb/instances/%i

# Optional: Initialize instance directory if it doesn't exist
# ExecStartPre=/opt/rustydb/current/bin/rusty-db-cli instance init --name %i --home /var/lib/rustydb/instances/%i

# Restart policy
Restart=on-failure
RestartSec=2
TimeoutStartSec=60
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true

# Allow write access only to instance directory
ReadWritePaths=/var/lib/rustydb/instances/%i

# File descriptor limits (high for database workloads)
LimitNOFILE=1048576
LimitNPROC=65536
LimitCORE=infinity

# Memory limits (optional, uncomment to set)
# MemoryMax=4G
# MemoryHigh=3G

# CPU limits (optional, uncomment to set)
# CPUQuota=200%

# OOM handling
OOMScoreAdjust=-500

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rustydb-%i

[Install]
WantedBy=multi-user.target
